# rc.local
#!/usr/bin/env sh
#
#

cd /root/Firewall

./Avvia.sh
